# package/secubox/luci-app-bandwidth-manager/root/etc/cron.d/bandwidth-alerts
# SecuBox Bandwidth Manager - Alert Monitoring
# Check bandwidth thresholds every 5 minutes
*/5 * * * * root /bin/sh -c '. /usr/share/secubox/lib/alerts.sh && run_alert_checks' >/dev/null 2>&1
