# package/secubox/secubox-core/root/etc/cron.d/secubox-integrity
# SecuBox File Integrity Monitor
# Check critical files every 5 minutes
*/5 * * * * root /usr/sbin/secubox-integrity check >/dev/null 2>&1
