# package/secubox/luci-app-crowdsec-dashboard/root/etc/cron.d/crowdsec-dashboard
# CrowdSec Dashboard cache refresh
# Refresh overview stats every minute - runs in background to avoid blocking rpcd
* * * * * root ( /usr/libexec/rpcd/luci.crowdsec-dashboard call refresh_cache & ) >/dev/null 2>&1
