# package/secubox/secubox-app-mitmproxy/files/etc/cron.d/waf-stats
# WAF stats cache updater
# Pre-generates stats JSON for fast LuCI dashboard loading
* * * * * root [ -x /usr/sbin/waf-stats-update ] && /usr/sbin/waf-stats-update >/dev/null 2>&1
