# package/secubox/secubox-vortex-dns/files/etc/init.d/vortex-dns
#!/bin/sh /etc/rc.common
# SecuBox Vortex DNS Daemon
# Meshed multi-dynamic subdomain delegation

START=95
STOP=10
USE_PROCD=1

PROG="/usr/sbin/vortexctl"

start_service() {
	local enabled
	config_load vortex-dns
	config_get enabled main enabled 0

	[ "$enabled" = "1" ] || return 0

	procd_open_instance
	procd_set_param command "$PROG" daemon
	procd_set_param respawn
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param pidfile /var/run/vortex-dns.pid
	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger "vortex-dns"
}

reload_service() {
	stop
	start
}
