# package/secubox/luci-app-secubox-netifyd/root/etc/crontabs/root.netifyd
# Netifyd flow data collector - runs every minute
* * * * * /usr/bin/netifyd-collector
