# package/secubox/secubox-app-watchdog/files/etc/cron.d/watchdog
# SecuBox Watchdog - Health check cron job
# Runs every minute to check and auto-recover services
* * * * * root /usr/sbin/watchdogctl check-recover >/dev/null 2>&1
