# package/secubox/luci-app-network-tweaks/root/etc/hotplug.d/iface/99-network-tweaks
#!/bin/sh
# Trigger network tweaks sync when LAN interface comes up

[ "$ACTION" = "ifup" ] || exit 0
[ "$INTERFACE" = "lan" ] || exit 0

# Wait a bit for network to stabilize
sleep 2

# Trigger sync
/usr/sbin/network-tweaks-sync sync >/dev/null 2>&1 &
