# package/secubox/luci-app-secubox-netifyd/root/etc/cron.d/netifyd-collector
# Netifyd Data Collector - runs every minute
* * * * * root /usr/sbin/netifyd-collector
